# deploy/systemd/nova-engine.service
[Unit]
Description=Nova Engine
After=network-online.target
Wants=network-online.target

[Service]
User=nova
WorkingDirectory=${NOVA_HOME}
EnvironmentFile=-${NOVA_HOME}/.env
ExecStart=${NOVA_HOME}/.venv/bin/python -u bot_main.py run
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
